# pkgs/voicevox-engine-cuda/voicevox-engine.service
[Unit]
Description=VOICEVOX engine HTTP server (CUDA build)
Documentation=https://github.com/VOICEVOX/voicevox_engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=voicevox-engine
Group=voicevox-engine
# Default: localhost のみ。 build host で
#   /etc/systemd/system/voicevox-engine.service.d/override.conf
# を用意して VV_HOST を tailscale0 IP に振る (ansible-nekonodesk 側の責務)。
Environment=VV_HOST=127.0.0.1
Environment=VV_PORT=50021
# bundle 同梱の CUDA 11 + cuDNN 8 + libpython3.11 stack を ld.so に見せる。
# 詳細は PKGBUILD の depends comment 参照。
Environment=LD_LIBRARY_PATH=/usr/lib/VOICEVOX/vv-engine
ExecStart=/usr/lib/VOICEVOX/vv-engine/run \
    --host=${VV_HOST} \
    --port=${VV_PORT} \
    --use_gpu \
    --voicevox_dir=/usr/lib/VOICEVOX/vv-engine \
    --runtime_dir=/usr/lib/VOICEVOX/vv-engine

StateDirectory=voicevox-engine
WorkingDirectory=%S/voicevox-engine
Restart=on-failure
RestartSec=5

# Hardening (= API server で reasonable な範囲)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# GPU device (/dev/nvidia*) を見せる必要があるため PrivateDevices=false
PrivateDevices=false
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-uvm-tools rw
DeviceAllow=/dev/nvidia-modeset rw
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
# onnxruntime/CUDA JIT で W+X mapping が必要
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
